# click_yes.service
[Unit]
Description=Click Yes

[Service]
Type=simple
ExecStart=/usr/local/sbin/click_yes.py

[Install]
WantedBy=multi-user.target
